# sys_files/usr/lib/systemd/system/buckle.socket
[Socket]
ListenStream=/trunk/socket/buckled.sock
Accept=no
RequiredBy=gild.service
